# digital-claude/assets/scripts/drain.service
[Unit]
Description=Background drain: manutenzione automatica skill digital-claude
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Sostituisci <PROJECT_PATH> con il path assoluto del progetto
# Es: /root/Progetti/digital-claude-dashboard
Environment=PROJECT_PATH=<PROJECT_PATH>
Environment=SKILL_DASHBOARD_URL=http://localhost:3001

ExecStart=/usr/bin/python3 \
  /root/.claude/skills/digital-claude/scripts/drain.py \
  --project ${PROJECT_PATH} \
  --curriculum-weekly

# Log su journal (visibile con: journalctl --user -u drain.service)
StandardOutput=journal
StandardError=journal

# Non fallire il timer se il progetto e' temporaneamente non disponibile
SuccessExitStatus=0

[Install]
WantedBy=default.target
